# deploy/systemd/atr-kraken.service
[Unit]
Description=ATR Kraken Engine Service
After=network.target

[Service]
Type=simple
User=dh
WorkingDirectory=/home/dh/serving-atr-inference
ExecStart=/home/dh/serving-atr-inference/.venvs/kraken/bin/python -m uvicorn kraken_svc.app:app --host 127.0.0.1 --port 8201
Restart=on-failure
RestartSec=5
Environment="CUDA_VISIBLE_DEVICES=1"

[Install]
WantedBy=multi-user.target