# deploy/systemd/atr-gateway.service
# asterAIx has no passwordless sudo and Linger=no, so prefer a USER unit:
#   mkdir -p ~/.config/systemd/user
#   cp deploy/systemd/atr-gateway.service ~/.config/systemd/user/
#   # one-time, needs an admin (survives logout):  sudo loginctl enable-linger tobias
#   systemctl --user daemon-reload && systemctl --user enable --now atr-gateway
#
# (If root is available, drop the same file in /etc/systemd/system and set User=tobias.)
[Unit]
Description=ATR serving gateway (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Repo/serving-atr-inference
EnvironmentFile=%h/Repo/serving-atr-inference/.env
ExecStart=%h/Repo/serving-atr-inference/.venvs/gateway/bin/uvicorn \
    atr_serving.app:app --host 0.0.0.0 --port 8200
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
