# deploy/systemd/atr-party.service
# User unit (systemctl --user) — asterAIx has no passwordless sudo / Linger=no.
# Install: scripts/install_user_units.sh  (copies to ~/.config/systemd/user/)
[Unit]
Description=ATR Party Engine Service (always-on)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Repo/serving-atr-inference/engines
EnvironmentFile=%h/Repo/serving-atr-inference/.env
Environment=CUDA_VISIBLE_DEVICES=1
Environment=PYTHONPATH=%h/Repo/serving-atr-inference/src
ExecStart=%h/Repo/serving-atr-inference/.venvs/party/bin/python -m uvicorn party_svc.app:app --host 127.0.0.1 --port 8203
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
